# infra/systemd/nexe-presence.service
[Unit]
Description=Nexe Presence Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/nexe/services/presence
EnvironmentFile=/var/www/html/nexe/services/presence/.env
ExecStart=/var/www/html/nexe/services/presence/bin/presence
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexe-presence

[Install]
WantedBy=multi-user.target
