# infra/systemd/nexe-notifications.service
[Unit]
Description=Nexe Notifications Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/html/nexe/services/notifications
EnvironmentFile=/var/www/html/nexe/services/notifications/.env
ExecStart=/var/www/html/nexe/services/notifications/bin/notifications
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexe-notifications

[Install]
WantedBy=multi-user.target
